# rpm/preinstall
#!/bin/bash

set +e
echo "Stopping Nginx"

systemctl stop nginx.service
systemctl disable nginx.service

exit 0
